# Config/Postgres/etc/profile.d/postgres_defaults.sh
#!/usr/bin/env bash
export PGDATA=/var/lib/pgsql/9.6/data

## Connection Defaults ##
#export PGHOST=localhost
export PGDATABASE=pardi
#export PGUSER=pardi_admin