# Config/Postgres/etc/profile.d/postgres_defaults.sh
#!/usr/bin/env bash
export PGDATA=/var/lib/pgsql/12/data

## Connection Defaults: use sockets and the logged in user name to connect ##
export PGDATABASE=ernie