# Config/Python/etc/profile.d/python_env.sh
#!/usr/bin/env bash
export PYTHONUNBUFFERED=1