# Config/Neo4j/etc/profile.d/env_neo4j.sh
#!/usr/bin/env bash
export NEO4J_HOME=/var/lib/neo4j